# agent/dist/fleet-agent.service
[Unit]
Description=Fleet Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fleet-agent run
Restart=always
RestartSec=5
KillMode=process

[Install]
WantedBy=multi-user.target
